# src/etc/cron.d/mrtg
0-55/5 *        * * *   root    if [ -x /usr/share/se3/sbin/mrtg-se3 ] && [ -r /etc/mrtg.cfg ]; then env LANG=C /usr/bin/mrtg /etc/mrtg.cfg >> /var/log/mrtg/mrtg.log 2>&1; env LANG=C /usr/share/se3/sbin/mrtg-se3 /etc/mrtg.cfg >> /var/log/mrtg/mrtg.log 2>&1; fi
