# sources/debian/sambaedu-boot-server.cron.d
*/5 * * * * root /usr/share/sambaedu/sbin/makedhcpconf.sh >/dev/null 2>&1
